# supervisor_conf/redis.conf
[program:Redis]
command=redis-server
directory=/tmp
stdout_logfile=%(ENV_HOME)s/log/youDL/redis.log
stderr_logfile=%(ENV_HOME)s/log/youDL/redis_errors.log
