# supervisor_conf/redis.conf
[program:Redis]
command=redis-server
directory=%(ENV_YOUDL_PATH)s/db
stdout_logfile=%(ENV_YOUDL_PATH)s/log/redis.log
stderr_logfile=%(ENV_YOUDL_PATH)s/log/redis_errors.log
